# powerlimit.service
[Unit]
Description=Power limit setter

[Service]
ExecStart=/home/s/powerlimit
Type=simple
Restart=no

[Install]
WantedBy=multi-user.target
